# docker/supervisor/cua-supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
user=root
loglevel=info
# Startup Performance Optimizations
minfds=4096
minprocs=512
identifier=cua-supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ==========================================
# Original Bytebot Services (Enhanced for C/ua)
# ==========================================

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24 -ac -nolisten tcp -dpi 96 +extension GLX +render -noreset
user=root
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb_error.log
environment=DISPLAY=":0"
# Performance Optimizations
startsecs=2
startretries=2
stopwaitsecs=5

[program:fluxbox]
command=/usr/bin/fluxbox
user=user
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox_error.log
environment=DISPLAY=":0",HOME="/home/user"
depends_on=xvfb
# Performance Optimizations
startsecs=3
startretries=2
stopwaitsecs=3

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 9990 --web /usr/share/novnc
user=user
autostart=true
autorestart=true
priority=150
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_error.log
environment=HOME="/home/user"
# Performance Optimizations - Start in parallel with x11vnc
startsecs=2
startretries=2
stopwaitsecs=3

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -nopw -listen localhost -xkb -ncache 10 -ncache_cr -noxrecord -noxfixes -noxdamage -wait 2 -shared -forever
user=root
autostart=true
autorestart=true
priority=150
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc_error.log
environment=DISPLAY=":0"
depends_on=xvfb
# Performance Optimizations - Reduced wait time and parallel startup
startsecs=3
startretries=2
stopwaitsecs=3

[program:bytebotd]
command=node /app/dist/main.js
directory=/app
user=user
autostart=true
autorestart=true
priority=300
stdout_logfile=/var/log/supervisor/bytebotd.log
stderr_logfile=/var/log/supervisor/bytebotd_error.log
environment=DISPLAY=":0",NODE_ENV="production",PORT="9990",NODE_OPTIONS="--max-old-space-size=512"
depends_on=fluxbox
# Performance Optimizations - Reduced dependencies and memory allocation
startsecs=5
startretries=2
stopwaitsecs=5

# ==========================================
# C/ua Framework Services
# ==========================================

[program:cua_agent_api]
command=/opt/cua/venv/bin/python -m uvicorn cua_agent.main:app --host 0.0.0.0 --port 9993 --workers 1 --timeout-keep-alive 5
directory=/opt/cua
user=user
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/cua_agent_api.log
stderr_logfile=/var/log/supervisor/cua_agent_api_error.log
environment=PYTHONPATH="/opt/cua",CUA_CONFIG_FILE="/opt/cua/config/agent-config.json",PYTHONUNBUFFERED="1"
# Performance Optimizations - Parallel startup with core services
startsecs=8
startretries=2
stopwaitsecs=5

[program:cua_websocket_server]
command=/opt/cua/venv/bin/python /opt/cua/scripts/websocket_server.py --host 0.0.0.0 --port 9996
directory=/opt/cua
user=user
autostart=true
autorestart=true
priority=250
stdout_logfile=/var/log/supervisor/cua_websocket.log
stderr_logfile=/var/log/supervisor/cua_websocket_error.log
environment=PYTHONPATH="/opt/cua",PYTHONUNBUFFERED="1"
# Performance Optimizations - Parallel startup
startsecs=5
startretries=2
stopwaitsecs=3

[program:cua_performance_monitor]
command=/opt/cua/venv/bin/python /opt/monitoring/scripts/performance_monitor.py --config /opt/cua/config/agent-config.json --port 9995
directory=/opt/monitoring
user=user
autostart=true
autorestart=true
priority=350
stdout_logfile=/var/log/supervisor/cua_monitor.log
stderr_logfile=/var/log/supervisor/cua_monitor_error.log
environment=PYTHONPATH="/opt/cua:/opt/monitoring",PYTHONUNBUFFERED="1"
# Performance Optimizations - Delayed start for performance monitoring
autostart=false
startretries=2
startsecs=5

[program:cua_health_checker]
command=/bin/bash -c 'while true; do /opt/cua/scripts/health-check.sh; sleep 30; done'
directory=/opt/cua
user=user
autostart=false
autorestart=true
priority=450
stdout_logfile=/var/log/supervisor/cua_health.log
stderr_logfile=/var/log/supervisor/cua_health_error.log
# Performance Optimizations - Start after core services and more frequent checks
startsecs=2
startretries=1

# ==========================================
# ANE Bridge Communication Monitor
# ==========================================

[program:ane_bridge_monitor]
command=/opt/cua/venv/bin/python /opt/ane-bridge/monitor.py
directory=/opt/ane-bridge
user=user
autostart=false
autorestart=true
priority=400
stdout_logfile=/var/log/supervisor/ane_monitor.log
stderr_logfile=/var/log/supervisor/ane_monitor_error.log
environment=PYTHONPATH="/opt/cua",ANE_BRIDGE_CONFIG="/opt/ane-bridge/bridge-config.json",PYTHONUNBUFFERED="1"
# Performance Optimizations - Optional service, start after core services are ready
startretries=2
startsecs=3

# ==========================================
# System Maintenance Services  
# ==========================================

[program:log_rotator]
command=/bin/bash -c 'while true; do find /var/log/supervisor -name "*.log" -size +100M -exec truncate -s 50M {} \; ; sleep 3600; done'
user=root
autostart=true
autorestart=true
priority=950
stdout_logfile=/dev/null
stderr_logfile=/var/log/supervisor/log_rotator_error.log

[program:cache_cleaner]
command=/bin/bash -c 'while true; do find /opt/cua/cache /opt/ane-bridge/cache -type f -mtime +1 -delete 2>/dev/null || true; sleep 1800; done'
user=user
autostart=true
autorestart=true
priority=960
stdout_logfile=/dev/null
stderr_logfile=/var/log/supervisor/cache_cleaner_error.log

# ==========================================
# Group Definitions for Optimized Parallel Startup
# ==========================================

# Primary services - Start first in parallel
[group:core_display]
programs=xvfb
priority=100

# Secondary services - Start in parallel after display is ready
[group:ui_services] 
programs=fluxbox,novnc,x11vnc
priority=150

# Application services - Start in parallel after UI is ready
[group:app_services]
programs=bytebotd,cua_agent_api,cua_websocket_server
priority=200

# Monitoring services - Start after applications are stable
[group:monitoring_services]
programs=cua_performance_monitor,cua_health_checker,ane_bridge_monitor
priority=350

# Maintenance services - Start last
[group:maintenance_services]
programs=log_rotator,cache_cleaner
priority=450

# ==========================================
# Service Orchestration Commands
# ==========================================

# Fast startup command for critical services only
[program:fast_startup]
command=/opt/cua/scripts/fast-startup.sh
user=user
autostart=false
autorestart=false
priority=50
stdout_logfile=/var/log/supervisor/fast_startup.log
stderr_logfile=/var/log/supervisor/fast_startup_error.log

# Delayed service starter for optional services
[program:delayed_services]
command=/bin/bash -c 'sleep 15; supervisorctl start cua_performance_monitor ane_bridge_monitor cua_health_checker'
user=user
autostart=true
autorestart=false
priority=500
stdout_logfile=/var/log/supervisor/delayed_services.log
stderr_logfile=/var/log/supervisor/delayed_services_error.log
startsecs=1